# training-club-overlay/skills/weekly-content-batch/weekly-content-batch.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.[USER].[AI_NAME].weekly-content-batch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd "$HOME/Documents/[AI_NAME]" &amp;&amp; /opt/homebrew/bin/claude -p "Run the weekly-content-batch skill — draft 7 days of social content for next week." 2&gt;&gt; "$HOME/Documents/[AI_NAME]/logs/weekly-content-batch.err"</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/[USER]</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/[USER]/Documents/[AI_NAME]/logs/weekly-content-batch.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/[USER]/Documents/[AI_NAME]/logs/weekly-content-batch.err</string>
</dict>
</plist>
